# tools/picoforge/deploy/service/probe/run
#!/bin/sh
# runit service: the availability probe.
#
# It runs forever in a 1-second loop, each iteration sending messages to the
# picomesh node to check whether the mesh's services are available: it polls
# the gateway's /_describe surface (which enumerates every active service) and
# the webapp page tier, and prints a one-line-per-service report to stdout
# (runsvdir → /init → /dev/console under the VM). runit keeps it supervised.
# PF_ROOT is normally exported by /opt/picoforge/init (and by the host run.sh);
# fall back to resolving it from $0's physical path for a manual launch.
PF_ROOT="${PF_ROOT:-$(cd "$(dirname "$(readlink -f "$0")")/../.." && pwd)}"
export PF_ROOT
. "${PF_ROOT}/service/common.sh"
cd "${PF_ROOT}"

# Services the collocated mesh node is expected to expose (the `plugins:` set
# from picoforge.yaml). The probe reports each as present/absent in /_describe.
EXPECTED_SERVICES="portalloc sharded_storage session accounts password_authn github_authn token_issuer issues git_repo git_pipeline personal_access_tokens relational_storage"

# Raw cookie-aware HTTP over busybox nc (the only client in the alpine guest
# that lets us read Set-Cookie). The server closes on `Connection: close`, so
# each call returns as soon as the response lands.
#   http_call PORT METHOD PATH BODY [extra-header] → raw response on stdout
http_call() {
    hc_port=$1; hc_method=$2; hc_path=$3; hc_body=$4; hc_hdr=${5:-}
    {
        printf '%s %s HTTP/1.1\r\n' "$hc_method" "$hc_path"
        printf 'Host: picomesh-yemu\r\n'
        [ -n "$hc_hdr" ] && printf '%s\r\n' "$hc_hdr"
        printf 'Content-Type: application/x-www-form-urlencoded\r\n'
        printf 'Content-Length: %s\r\n' "${#hc_body}"
        printf 'Connection: close\r\n\r\n'
        printf '%s' "$hc_body"
    } | nc -w 30 127.0.0.1 "$hc_port" 2>/dev/null
}
sid_of() { tr -d '\r' | sed -n 's/.*[Ss]et-[Cc]ookie:[ ]*picomesh-sid=\([^;]*\).*/\1/p' | head -1; }

# Seed a usable demo on a fresh boot, idempotent (duplicates just come back
# "already taken"):
#   - root/rootpw  — the FIRST account registered becomes the site owner/admin.
#   - demo/demo    — a regular user...
#   - demo/hello, demo/world — ...with two repos, created through the WEBAPP the
#     same way a browser does (POST /-/repos/new with the session cookie). This
#     is what makes the browser demo land usable and what the headless-chrome
#     test asserts (login demo/demo → /demo lists hello + world).
seed() {
    command -v nc >/dev/null 2>&1 || { say probe "no nc — cannot seed"; return; }
    wait_http "http://127.0.0.1:${MESH_PORT}/_describe" 120 \
        || { say probe "gateway never came up — skipping seed"; return; }

    # 1. site owner (must be first).
    http_call "${MESH_PORT}" POST /register 'username=root&password=rootpw' \
        | head -1 | grep -q '30[0-9]' \
        && say probe "seeded site-owner root/rootpw" \
        || say probe "root/rootpw already present"

    # 2. demo user — capture its session id (register, else log in).
    sid=$(http_call "${MESH_PORT}" POST /register 'username=demo&password=demo' | sid_of)
    [ -z "$sid" ] && sid=$(http_call "${MESH_PORT}" POST /login 'username=demo&password=demo' | sid_of)
    if [ -z "$sid" ]; then say probe "WARN: no demo session — skipping repo seed"; return; fi
    say probe "seeded demo/demo"

    # 3. demo's repos, via the webapp (browser-equivalent). Wait for it first.
    wait_http "http://127.0.0.1:${WEBAPP_HTTP_PORT}/-/login" 60 \
        || say probe "webapp not up — skipping repo seed"
    for repo in hello world; do
        loc=$(http_call "${WEBAPP_HTTP_PORT}" POST /-/repos/new "name=${repo}" \
                "Cookie: picomesh-sid=${sid}" \
              | tr -d '\r' | sed -n 's/^[Ll]ocation:[ ]*//p' | head -1)
        case "$loc" in
            */demo/${repo}*) say probe "created repo demo/${repo}" ;;
            *)               say probe "repo demo/${repo} -> ${loc:-(exists?)}" ;;
        esac
    done
}
seed

# Forever: poll the mesh + webapp, print availability to stdout, sleep 1s.
while :; do
    describe=$(wget -q -O - -T 2 "http://127.0.0.1:${MESH_PORT}/_describe" 2>/dev/null)
    if [ -z "${describe}" ]; then
        say probe "gateway :${MESH_PORT} unreachable"
    else
        up=0; down=0
        for svc in ${EXPECTED_SERVICES}; do
            if printf '%s' "${describe}" | grep -q "\"service\":\"${svc}\""; then
                say probe "available: ${svc}"
                up=$((up + 1))
            else
                say probe "MISSING:   ${svc}"
                down=$((down + 1))
            fi
        done
        say probe "summary: ${up} up, ${down} missing"
    fi

    # Probe the real login page (/-/login — a command path, served without auth
    # and with NO backend calls). A bare path like /login is a NAMESPACE lookup
    # that would drive git_repo.*_for_namespace on every poll.
    if wget -q -O /dev/null -T 2 "http://127.0.0.1:${WEBAPP_HTTP_PORT}/-/login" 2>/dev/null; then
        say probe "available: webapp page tier :${WEBAPP_HTTP_PORT}"
    else
        say probe "MISSING:   webapp page tier :${WEBAPP_HTTP_PORT}"
    fi

    sleep 1
done
